# politeping/deploy/politeping.service
[Unit]
Description=PolitePing FastAPI
After=network-online.target
Wants=network-online.target

[Service]
User=www-data
WorkingDirectory=/opt/politeping
EnvironmentFile=/opt/politeping/.env
ExecStart=/opt/politeping/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --proxy-headers --workers 1
Restart=always
RestartSec=3
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target